# Friendlywrt_archive/check
#!/bin/sh /etc/rc.common

START=99
 
start() {
    nohup /bin/sh /usr/bin/check_inet.sh 1>/dev/null 2>&1 &
}